# files/root/etc/supervisor/conf.d/postgres.conf
[program:postgresql]
command=/usr/lib/postgresql/9.1/bin/postgres
user=postgres
autorestart=true
