# files/root/etc/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=/usr/lib/postgresql/9.5/bin/postgres
user=postgres
autorestart=true
